# zfiso/runit/admin/run
#!/bin/bash -e
exec 2>&1
cd /usr/local/admin/
exec chpst -u root ./admin 
