# ek/runit/elasticsearch/run
#!/bin/bash -e
exec 2>&1
cd /usr/share/elasticsearch
exec chpst -u elasticsearch bin/elasticsearch -p pid
