# zfiso/runit/logstash/run
#!/bin/bash -e
exec 2>&1
cd /usr/share/logstash/
exec chpst -u root bin/logstash -f /etc/logstash/conf.d/ 
