# zfiso/runit/elasticsearch/run
#!/bin/bash -e
exec 2>&1
cd /usr/local/elasticsearch-2.4.3
exec chpst -u elasticsearch bin/elasticsearch 
