# zfiso/runit/authority/run
#!/bin/sh 
exec 2>&1
cd /usr/local/authority
exec chpst -u root ./authority 
